# backend/copytrader-backend.service
[Unit]
Description=CopyTradinj Backend Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/copytradinj/backend
Environment=NODE_ENV=production
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=copytradinj-backend

[Install]
WantedBy=multi-user.target
